# templates/default/gnrd_worker.conf.erb
description 'start and stop resque workers for GNRD'
version '0.01'
author 'Dmitry Mozzherin'

start on runlevel [2345]

stop on runlevel [016]

respawn

expect fork

script
chdir <%= node['gnrd']['directory'] %>
exec RACK_ENV=production COUNT=<%= node['gnrd']['worker_count'] %> QUEUE=name_finder rake resque:workers  &
emit gnrd_worker_starting
end script
